# vpnc_client/rootfs/etc/services.d/vpnc_client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

## Get the keys from the user config options.
gateway=$(bashio::config "vpn_server.gateway")
ID=$(bashio::config "vpn_server.ID")
secret=$(bashio::config "vpn_server.secret")
username=$(bashio::config "vpn_server.username")
password=$(bashio::config "vpn_server.password")
server=$(bashio::config 'ping.server')
wait=$(bashio::config 'ping.wait' 5)
debug=$(bashio::config 'debug' 0)

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "Creating vpnc config file"
VPNC_CONFIG_PATH="/etc/vpnc/connection.conf"
cat > "$VPNC_CONFIG_PATH" <<EOL
IPSec gateway $gateway
IPSec ID $ID
IPSec secret $secret
IKE Authmode psk
Xauth username $username
Xauth password $password
Script /etc/vpnc/vpnc-script-no-default-route
local port 0
DPD idle timeout (our side) 0
Debug $debug
EOL


## Run your program
bashio::log.info "Starting vpnc connection"
vpnc connection

## Add custom routes through VPN
if bashio::config.has_value 'vpn_routes'; then
    bashio::log.info "Adding custom VPN routes"
    for route in $(bashio::config 'vpn_routes'); do
        bashio::log.info "Adding route: $route via tun0"
        ip route add "$route" dev tun0 2>/dev/null || bashio::log.warning "Route $route already exists or failed"
    done
fi

sleep $wait

while ping -c 1 $server &> /dev/null
do
  bashio::log.info "Ping okay, wait ${wait}s"
  sleep $wait
done

bashio::log.error "Connection lost"
