# vpnc_client/rootfs/etc/services.d/vpnc_client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
# declare gateway
# declare ID
# declare secret
# declare username
# declare password
# declare server
# declare wait
# declare debug  

## Get the keys from the user config options.
gateway=$(bashio::config "vpn_server.gateway")
ID=$(bashio::config "vpn_server.ID")
secret=$(bashio::config "vpn_server.secret")
username=$(bashio::config "vpn_server.username")
password=$(bashio::config "vpn_server.password")
server=$(bashio::config 'ping.server')
wait=$(bashio::config 'ping.wait' 5)
debug=$(bashio::config 'debug' 0)

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "creating vpnc config file"
VPNC_CONFIG_PATH="/etc/vpnc/connection.conf"
cat > "$VPNC_CONFIG_PATH" <<EOL
IPSec gateway $gateway
IPSec ID $ID
IPSec secret $secret
IKE Authmode psk
Xauth username $username
Xauth password $password
local port 0
DPD idle timeout (our side) 0
Debug $debug
EOL


## Run your program
bashio::log.info "starting vpnc connection"
vpnc connection
sleep $wait

while ping -c 1 $server &> /dev/null
do
  bashio::log.info "ping okay, wait ${wait}s"
  sleep $wait
done

